# deploy/supervisor/include/exporter_20001000005.conf

[program:exporter_20001000005]
command=python -m apps.exporter.script.exporter
environment=MODULE_ID="20001000005",MODULE_SECRET="123456",MODULE_URL="http://www.yunplc.com:7080",RANDOM_PORT=20000,HOST="0.0.0.0",ADVERTISE="host.docker.internal"
autostart=false

[module]
name=收集器
module_id=20001000005
module_secret=123456
module_url=http://www.yunplc.com:7080
export_interval=5
export_timeout=3
